# systemd/firewallbot.service
[#] FireWallBot systemd unit (generated by install script or manual edit)
[Unit]
Description=FireWallBot Login Watcher
After=network.target

[Service]
Type=simple
# WorkingDirectory and Environment should be set by install script to your repo path.
WorkingDirectory=@REPO@
Environment=FIREWALLBOT_LOG_DIR=@REPO@/log
ExecStart=/usr/bin/env python3 scripts/login_watcher.py
Restart=always
RestartSec=2s
User=root
Group=root

[Install]
WantedBy=multi-user.target

