# scripts/LoginWatcher/LoginWatcher.service.tmpl
[Unit]
Description=FireWallBot LoginWatcher
After=network.target

[Service]
Type=simple
WorkingDirectory=@REPO@
Environment=FIREWALLBOT_LOG_DIR=@REPO@/log
ExecStart=/usr/bin/env python3 scripts/LoginWatcher/watcher.py
Restart=always
RestartSec=2s
User=root
Group=root

[Install]
WantedBy=multi-user.target

