# scripts/loginwatcher/loginwatcher.service.tmpl
[Unit]
Description=FireWallBot loginwatcher
After=network.target

[Service]
Type=simple
WorkingDirectory=@REPO@
Environment=FIREWALLBOT_LOG_DIR=@REPO@/log
ExecStart=/usr/bin/env python3 scripts/loginwatcher/watcher.py
Restart=always
RestartSec=2s
User=root
Group=root

[Install]
WantedBy=multi-user.target
